# build/rpstir2.service
[Unit]
Description=rpstir2 service
After=mysqld.service
#Requires=mysqld.service
[Service]
Type=forking
WorkingDirectory=/root/rpki/rpstir2/bin
ExecStart=/root/rpki/rpstir2/bin/rpstir2-service.sh start
#Restart=always
[Install]
WantedBy=multi-user.target